# dotfile-watcher/.config/systemd/user/dotfile-watcher.path
[Unit]
Description=Watch dotfiles directory for changes

[Path]
PathChanged=%h/dotfiles/.env
PathChanged=%h/dotfiles/envs
PathChanged=%h/dotfiles/git

[Install]
WantedBy=default.target
