# dotfile-watcher/.config/systemd/user/dotfile-watcher.service
[Unit]
Description=Dotfile watcher

[Service]
ExecStart=%h/.local/bin/watch-dotfiles.sh
Restart=always

[Install]
WantedBy=default.target
