# dotfile-watcher/.config/systemd/user/dotfile-watcher.service
[Unit]
Description=Regenerate templates

[Service]
Type=oneshot
WorkingDirectory=%h/dotfiles
ExecStart=/usr/bin/make templates
